# crontab
0 * * * * /home/ubuntu/bin/vault-snapshot-s3.sh
